# config/com.example.whispervoicecontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.whispervoicecontrol</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/martinstuessy/Code/whisper-yabai-mac-os-x/venv/bin/python</string>
        <string>/Users/martinstuessy/Code/whisper-yabai-mac-os-x/src/simplified_daemon.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/martinstuessy/Code/whisper-yabai-mac-os-x</string>
    <key>StandardErrorPath</key>
    <string>/Users/martinstuessy/Code/whisper-yabai-mac-os-x/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/martinstuessy/Code/whisper-yabai-mac-os-x/output.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>